# src/basil_templates/basil_flask/project_base/puppet/modules/flask/files/etc/init/runflask.conf
# Basil: Build A System Instant-Like
# Copyright (C) 2014 Catalyst IT Ltd
#
# This program is free software: you can redistribute it and/or modify
# it under the terms of the GNU Affero General Public License as published by
# the Free Software Foundation, either version 3 of the License, or
# (at your option) any later version.
#
# This program is distributed in the hope that it will be useful,
# but WITHOUT ANY WARRANTY; without even the implied warranty of
# MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.  See the
# GNU General Public License for more details.
#
# You should have received a copy of the GNU Affero General Public License
# along with this program.  If not, see <http://www.gnu.org/licenses/>.

# runflask - runflask job file

description "Run Flask app"

# Stanzas
#
# Stanzas control when and how a process is started and stopped
# See a list of stanzas here: http://upstart.ubuntu.com/wiki/Stanzas#respawn

# When to start the service
start on vagrant-mounted

# When to stop the service
stop on runlevel [!2345]

# Automatically restart process if crashed
respawn

# Start the process
exec /usr/bin/python /home/vagrant/src/main.py